# deploy/ady-ticket-bot.service
[Unit]
Description=Ady.az Baku<->Tbilisi ticket watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ady-bot
WorkingDirectory=/opt/ady-ticket-checker-bot
EnvironmentFile=/opt/ady-ticket-checker-bot/.env
ExecStart=/opt/ady-ticket-checker-bot/.venv/bin/python -m ady_ticket_bot.main
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
